# ipmitool/services/ipmi-monitor.service
[Unit]
Description=IPMI Fan Speed and Sensor Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/ipmi_monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target